# debian/osm-tile-server.cron.daily
#!/bin/sh


if [ -e /var/lib/naturfakta-web/images ] ; then
  chgrp www-data `readlink -f /var/lib/naturfakta-web/images`
  chmod g+w `readlink -f /var/lib/naturfakta-web/images`
  mkdir -p /var/lib/naturfakta-web/images/orig
  chgrp www-data `readlink -f /var/lib/naturfakta-web/images/orig`
  chmod g+w `readlink -f /var/lib/naturfakta-web/images/orig`
  cp /usr/share/naturfakta-web/keys/check.php /usr/share/naturfakta-web/keys/.htaccess /var/lib/naturfakta-web/keys/
  ln -sf /usr/share/websites/naturfakta/www/versions/bimage.php /var/lib/naturfakta-web/images/bimage.php
fi
